# ct/ct-cubieez/overlay/etc/udev/rules.d/11-cubie.rules
DEVPATH=="/kernel/disp/hdmi", RUN+="/bin/disp.sh"
